# pi_pianoteq.service
[Unit]
Description = Service for pi-pianoteq

[Service]
Type=simple
ExecStart=python3 -m pi-pianoteq

[Install]
WantedBy=default.target
